# http-redirect/redirect-http.socket
[Unit]
Description=Redirect HTTP to HTTPS socket

[Socket]
ListenStream=80
Accept=yes

[Install]
WantedBy=sockets.target
