# http-redirect/redirect-http@.service
[Unit]
Description=Redirect HTTP to HTTPS Per-connection

[Service]
ExecStart=/usr/bin/awk -f /usr/bin/redirect.awk
StandardInput=socket
